# redis.conf
[program:redis]
command=redis-entrypoint.sh redis-server
